# scripts/com.user.claude-audit.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!--
    LAUNCHD PLIST TEMPLATE — Claude Code quarterly config audit.

    Before using:
    1. Replace `com.user.claude-audit` (in <key>Label</key> below) with your own
       reverse-DNS identifier, e.g. `com.yourdomain.claude-audit`.
    2. Replace every `${HOME}` placeholder with your absolute home directory path.
       launchd does NOT expand `$HOME` or `~` — paths must be absolute.
       Example: /Users/yourname/.claude/scripts/audit-config.sh
    3. Rename this file to match your Label (e.g. `com.yourdomain.claude-audit.plist`)
       and save it to `~/Library/LaunchAgents/`.
    4. Load it: `launchctl load ~/Library/LaunchAgents/com.yourdomain.claude-audit.plist`
    5. Verify: `launchctl list | grep claude-audit`
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.claude-audit</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>${HOME}/.claude/scripts/audit-config.sh</string>
    </array>

    <!--
        Quarterly schedule: first day of Jan / Apr / Jul / Oct at 03:00 local.
        StartCalendarInterval with multiple dicts = OR semantics — fires on
        any matching condition.
    -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Month</key><integer>1</integer>
            <key>Day</key><integer>1</integer>
            <key>Hour</key><integer>3</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Month</key><integer>4</integer>
            <key>Day</key><integer>1</integer>
            <key>Hour</key><integer>3</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Month</key><integer>7</integer>
            <key>Day</key><integer>1</integer>
            <key>Hour</key><integer>3</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Month</key><integer>10</integer>
            <key>Day</key><integer>1</integer>
            <key>Hour</key><integer>3</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
    </array>

    <!-- Run even if the laptop was asleep at the scheduled time -->
    <key>StartCalendarIntervalRunsWhileSleeping</key>
    <false/>

    <!-- Logs -->
    <key>StandardOutPath</key>
    <string>${HOME}/.claude/audit-reports/launchd.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>${HOME}/.claude/audit-reports/launchd.stderr.log</string>

    <!-- Don't restart if it fails — quarterly run, not a daemon -->
    <key>KeepAlive</key>
    <false/>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
